# findword.service
[Unit]
Description=Find log file
Wants=findword.timer

[Service]
Type=oneshot
ExecStart=/bin/grep "systemd" /var/log/messages
