# findword.timer
[Unit]
Description=Running findwork.service

[Timer]
OnUnitActiveSec=30
Unit=findword.service

[Install]
WantedBy=timers.target
