# findword.timer
[Unit]
Description=Running findwork.service

[Timer]
OnCalendar=*:0/30
Unit=findword.service

[Install]
WantedBy=multi-user.target
